# com.annchiahui.woo-sprinkles.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.annchiahui.woo-sprinkles.watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/anntai/tools/woo-sprinkles/watch.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/anntai</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/woo-sprinkles-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/woo-sprinkles-watch.err</string>
</dict>
</plist>
