# com.annchiahui.woo-sprinkles.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.annchiahui.woo-sprinkles.menubar</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/anntai/tools/woo-sprinkles/WooSprinklesMenuBar.app/Contents/MacOS/WooSprinklesMenuBar</string>
        <string>cyan</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/anntai</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/woo-sprinkles-menubar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/woo-sprinkles-menubar.err</string>
</dict>
</plist>
